# openclaw/prod/systemd/revenue-os-prod-scheduler.service
[Unit]
Description=Revenue OS scheduler prod
After=network-online.target

[Service]
Type=oneshot
User=revenueos
Group=revenueos
WorkingDirectory=/opt/revenue-os
Environment=NODE_ENV=production
Environment=REVENUE_OS_ENVIRONMENT=prod
Environment=NODE_COMPILE_CACHE=/var/tmp/openclaw-compile-cache
Environment=OPENCLAW_NO_RESPAWN=1
EnvironmentFile=-/opt/revenue-os/.secrets/revenue-os.local.env
ExecStart=/usr/bin/env bash -lc 'cd /opt/revenue-os && npm run runtime:scheduler'
