# openclaw/prod/systemd/revenue-os-prod.service
[Unit]
Description=Revenue OS OpenClaw prod gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/revenue-os
Environment=NODE_ENV=production
ExecStart=/usr/bin/env openclaw gateway --config %h/revenue-os/openclaw/prod/openclaw.json
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
