# openclaw/prod/systemd/revenue-os-prod-source-refresh.timer
[Unit]
Description=Revenue OS source-refresh timer prod

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
Unit=revenue-os-prod-source-refresh.service
Persistent=true

[Install]
WantedBy=timers.target
