# openclaw/prod/systemd/revenue-os-prod-source-refresh.service
[Unit]
Description=Revenue OS source-refresh prod
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/revenue-os
Environment=NODE_ENV=production
Environment=REVENUE_OS_ENVIRONMENT=prod
EnvironmentFile=-%h/revenue-os/.secrets/revenue-os.local.env
ExecStart=/usr/bin/env bash -lc 'cd %h/revenue-os && npm run refresh:updates'
