# openclaw/stage/systemd/revenue-os-stage.service
[Unit]
Description=Revenue OS OpenClaw stage gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=revenueos
Group=revenueos
WorkingDirectory=/opt/revenue-os
Environment=NODE_ENV=production
Environment=REVENUE_OS_ENVIRONMENT=stage
EnvironmentFile=-/opt/revenue-os/.secrets/revenue-os.local.env
ExecStartPre=/usr/bin/env bash -lc 'cd /opt/revenue-os && npm run runtime:probe-models && npm run bootstrap:control-plane'
ExecStartPre=/usr/bin/env bash -lc 'cd /opt/revenue-os && openclaw doctor'
ExecStart=/usr/bin/env openclaw gateway --config /opt/revenue-os/openclaw/stage/openclaw.json
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
