# openclaw/stage/systemd/revenue-os-stage.service
[Unit]
Description=Revenue OS OpenClaw stage gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=revenueos
Group=revenueos
WorkingDirectory=/opt/revenue-os
Environment=NODE_ENV=production
Environment=REVENUE_OS_ENVIRONMENT=stage
Environment=NODE_COMPILE_CACHE=/var/tmp/openclaw-compile-cache
Environment=OPENCLAW_NO_RESPAWN=1
EnvironmentFile=-/opt/revenue-os/.secrets/revenue-os.local.env
Environment=OPENCLAW_CONFIG_PATH=/opt/revenue-os/data/generated/openclaw/stage.json
ExecStartPre=/usr/bin/env bash -lc 'cd /opt/revenue-os && export OPENCLAW_BIN="${OPENCLAW_BIN:-$(/opt/revenue-os/scripts/bootstrap/resolve-openclaw-bin.sh /opt/revenue-os)}" && npm run runtime:probe-models -- --active && npm run bootstrap:control-plane'
ExecStartPre=/usr/bin/env bash -lc 'cd /opt/revenue-os && bash scripts/bootstrap/sync-runtime-workspace.sh stage'
ExecStartPre=/usr/bin/env bash -lc 'cd /opt/revenue-os && export OPENCLAW_BIN="${OPENCLAW_BIN:-$(/opt/revenue-os/scripts/bootstrap/resolve-openclaw-bin.sh /opt/revenue-os)}" && npm run runtime:render-openclaw-config -- --environment stage --output /opt/revenue-os/data/generated/openclaw/stage.json'
ExecStartPre=/usr/bin/env bash -lc 'cd /opt/revenue-os && export OPENCLAW_BIN="${OPENCLAW_BIN:-$(/opt/revenue-os/scripts/bootstrap/resolve-openclaw-bin.sh /opt/revenue-os)}" && npm run verify:openclaw-config -- stage'
ExecStartPre=/usr/bin/env bash -lc 'cd /opt/revenue-os && export OPENCLAW_BIN="${OPENCLAW_BIN:-$(/opt/revenue-os/scripts/bootstrap/resolve-openclaw-bin.sh /opt/revenue-os)}" && OPENCLAW_CONFIG_PATH=/opt/revenue-os/data/generated/openclaw/stage.json "$OPENCLAW_BIN" doctor'
ExecStart=/usr/bin/env bash -lc 'cd /opt/revenue-os && export OPENCLAW_BIN="${OPENCLAW_BIN:-$(/opt/revenue-os/scripts/bootstrap/resolve-openclaw-bin.sh /opt/revenue-os)}" && OPENCLAW_CONFIG_PATH=/opt/revenue-os/data/generated/openclaw/stage.json "$OPENCLAW_BIN" gateway'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
